# macsmc-battery/systemd/macsmc-battery-charge-control-end-threshold.path
[Unit]
Description=Monitor charge_control_end_threshold

[Path]
PathChanged=/sys/class/power_supply/macsmc-battery/charge_control_end_threshold
