# systemd/sd-asahi-early-firmware.target
[Unit]
Description=Asahi Early Firmware
Before=systemd-udevd.service
Requires=sd-asahi-early-firmware.service

[Install]
WantedBy=sysinit.target
