# systemd/sd-asahi-late-firmware.target
[Unit]
Description=Asahi Late Firmware
After=initrd-fs.target cryptsetup.target
Requires=sd-asahi-late-firmware.service

[Install]
WantedBy=initrd.target
